# containers/monitor/supervisord.conf
[supervisord]
childlogdir = /var/log
logfile = /var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 1
loglevel = info
pidfile = /var/run/supervisord.pid
nodaemon = true
nocleanup = false
environment = C_FORCE_ROOT=1

[inet_http_server]
port = 0.0.0.0:9001
username =
password =

[supervisorctl]
serverurl = http://localhost:9001
username =
password =

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:bucky]
command = python /mist.monitor/src/bucky/bucky.py /bucky_conf.py
process_name = bucky
priority = 50

[program:monitor]
command = uwsgi --plugin /usr/lib/uwsgi/python_plugin.so --http-socket 0.0.0.0:80 --paste-logger --ini-paste /uwsgi.ini
directory = /mist.monitor
process_name = monitor
priority = 60

[program:alert]
command = mist-alert
directory = /mist.monitor
process_name = alert
priority = 60
