# conf/supervisord.conf.in

[supervisord]
childlogdir = ${buildout:directory}/var/log
logfile = ${buildout:directory}/var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = ${buildout:directory}/var/supervisord.pid
umask = 022
nodaemon = false
nocleanup = false
environment = LD_LIBRARY_PATH=${buildout:directory}/parts/mongo-c-driver/lib,PERL5LIB=${buildout:directory}/parts/collectd/lib/perl5/:${buildout:directory}/parts/collectd-graphite/lib/perl5/

[inet_http_server]
port = 127.0.0.1:9010
username = mistiera
password =

[supervisorctl]
serverurl = http://127.0.0.1:9010
username = mistiera
password =

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:uwsgi]
command = ${buildout:directory}/bin/uwsgi --ini-paste ${buildout:directory}/production.ini -x ${buildout:directory}/parts/uwsgi/uwsgi.xml --paste-logger
process_name = uwsgi
directory = ${buildout:directory}
priority = 10
redirect_stderr = false
user = mist

[program:bucky]
command = ${buildout:directory}/pypy-env/bin/bucky ${buildout:directory}/conf/bucky_conf.py
process_name = bucky
directory = ${buildout:directory}
priority = 15
redirect_stderr = false
autorestart = true
user = mist

[program:mongodb]
command = ${buildout:directory}/bin/mongod  --journal --dbpath ${buildout:directory}/var/data --bind_ip 127.0.0.1 --master --port 27022 --nohttpinterface --logpath ${buildout:directory}/var/log/mongodb.log --directoryperdb
process_name = mongodb
directory = ${buildout:directory}/bin
priority = 30
redirect_stderr = false
user = mist

[program:alert]
command = ${buildout:directory}/pypy-env/bin/pypy ${buildout:directory}/src/mist/alert/alert.py
process_name = alert
directory = ${buildout:directory}
priority = 40
redirect_stderr = false
user = mist

[program:memcache]
command = ${buildout:directory}/bin/memcached  -l 127.0.0.1
process_name = memcache
directory = ${buildout:directory}/bin
priority = 25
redirect_stderr = false
user = mist

